# systems/lossless_player.service
[Unit]
Description=Lossless Audio Player for Raspberry Pi
After=network.target sound.target

[Service]
User=jay
Group=jay
WorkingDirectory=/home/jay/lossless_player
ExecStart=/home/jay/lossless_player/venv/bin/python3 /home/jay/lossless_player/src/main.py
Restart=always
RestartSec=5
StandardOutput=syslog
StandardError=syslog
SyslogIdentifier=lossless_player
Environment="XDG_RUNTIME_DIR=/run/user/1000"

[Install]
WantedBy=multi-user.target
